# 极验证/jiyanzheng.conf
[program:极验证]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000
directory=/root/极验证
autostart=true
autorestart=true
environment=LANG="en_US.utf8", LC_ALL="en_US.UTF-8", LC_LANG="en_US.UTF-8"